# root/etc/periodic/daily/clean_syslogs.sh
#!/usr/bin/env bash

find /var/log/ -mtime +30 -delete
